# systemd/override/anbox-session-manager.service
[Service]
ExecStart=
ExecStart=/usr/bin/anbox session-manager --server-side-decoration --no-touch-emulation
Restart=on-failure
RestartSec=3