# systemd/override/anbox-session-manager.service
[Service]
ExecStart=
ExecStart=/usr/bin/anbox session-manager --server-side-decoration
Restart=on-failure
RestartSec=3

[Install]
Alias=anbox.service