# systemd/override/anbox-session-manager.service
[Unit]
Before=anbox-launchers.service anbox-launchers.timer
Requires=anbox-launchers.service anbox-launchers.timer
