# systemd/override/epic-claim.timer
[Timer]
OnCalendar=
OnCalendar=18:00 UTC
