# systemd/timer/aur-autobuild.timer
[Unit]
Description=Build aur packages weekly.

[Timer]
OnCalendar=weekly
RandomizedDelaySec=4h
Persistent=true

[Install]
WantedBy=timers.target
