# systemd/timer/vpnbook-password.timer
[Unit]
Description=Update vpnbook password weekly.

[Timer]
OnCalendar=weekly
RandomizedDelaySec=1d
Persistent=true
[Install]
WantedBy=timers.target
