# systemd/service/github-sync-forks.service
[Unit]
Description=sync github forks

[Service]
MemoryMax=100M
ExecStart=%h/.config/systemd/scripts/github-sync-forks
