# systemd/service/trakt-scrobbler.service
[Unit]
Description=Trakt Scrobbler Service

[Service]
ExecStart="/usr/bin/trakts" run

[Install]
WantedBy=default.target