# systemd/service/vpnbook-password.service
[Unit]
Description=vpnbook password updater.
[Service]
ExecStart=/bin/env %h/bin/vpnbook update
[Install]
WantedBy=default.target
