# systemd/unit/clamscan.service
[Unit]
Description=clamav scanner

[Service]
MemoryHigh=100M
MemoryMax=300M
ExecStart=%h/.config/systemd/scripts/clamscan
