# systemd/unit/gpg-refresh.service
[Unit]
Description=Refresh local keyring

[Service]
ExecStart=/usr/bin/gpg --refresh-keys
SuccessExitStatus=2
