# systemd/unit/trakt-scrobbler.service
[Unit]
Description=Trakt Scrobbler Service

[Service]
ExecStart=python -m trakt_scrobbler run

[Install]
WantedBy=default.target
