# systemd/unit/vpnbook.service
[Unit]
Description=vpnbook.net password updater

[Service]
ExecStart=%h/.config/systemd/scripts/vpnbook
