# systemd/unit/vpnbook.timer
[Unit]
Description=Update vpnbook.net password daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
