# picom/dot-config/systemd/user/picom@.service
[Unit]
Description=PiCom Compositning Manager on Display %I
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=!XDG_SESSION_DESKTOP=KDE
ConditionEnvironment=!XDG_SESSION_DESKTOP=GNOME

[Service]
Environment=DISPLAY="%I"
ExecStart=/usr/bin/picom
ExecReload=kill -USR1 $MAINPID
RestartSec=3
Restart=on-failure
MemoryMax=512M
MemoryHigh=256M
CPUQuota=5%
MemoryAccounting=yes

[Install]
WantedBy=graphical-session.target